# deploy/nestor-matrix.service
[Unit]
Description=Néstor Matrix Bot
Documentation=https://github.com/elatomo/nestor-matrix
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nestor
Group=nestor
WorkingDirectory=/opt/nestor-matrix/repo
EnvironmentFile=/opt/nestor-matrix/repo/.env
ExecStart=/opt/nestor-matrix/.local/bin/uv run nestor-matrix run

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6
SystemCallArchitectures=native
ReadWritePaths=/opt/nestor-matrix

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
